# software/ubuntu_conf/ttymxc0.conf
start on stopped rc or RUNLEVEL=[12345]
stop on runlevel [!12345]

respawn
exec /sbin/getty -L 115200 console vt102
